# conf/saltstack/srv/salt/qa/conf.d/crontab/client.conf
# Lines below here are managed by Salt, do not edit
# SALT_CRON_IDENTIFIER:CRON-NTPDATE
*/20 * * * * /usr/sbin/ntpdate ntp.office.test >/dev/null 2>&1 &
